# deploy/systemd/stocktalk-stop.service
[Unit]
Description=Stop Stocktalk service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl stop stocktalk.service
