# deploy/systemd/stocktalk-start.service
[Unit]
Description=Start Stocktalk service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start stocktalk.service
